# supervisord.d/bitbns_live_prices_v2.ini
[program:bitbns_live_prices_v2]
command=/opt/scripts/bitbns_live_prices_v2.py
stderr_logfile=/var/log/bitbns_live_prices_v2.log
stdout_logfile=/var/log/bitbns_live_prices_v2_err.log
autostart=true
autorestart=true
logfile_maxbytes = 10MB
logfile_backups=2
startretries=10
