# supervisord.d/koinex_live_prices_v2.ini
[program:koinex_live_prices_v2]
command=/opt/scripts/koinex_live_prices_v2.py
stderr_logfile=/var/log/koinex_live_prices_v2.log
stdout_logfile=/var/log/koinex_live_prices_v2.log
autostart=true
autorestart=true
logfile_maxbytes = 10MB
logfile_backups=2
startretries=10
